# scripts/auto_resume.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  macOS only — this is a launchd agent. On Linux use a systemd timer or cron entry that runs
  the same auto_resume.sh; on Windows, a Task Scheduler task. The script is the portable part.

  Replace __HARNESS_ROOT__ with your checkout's absolute path before loading:
    sed -i '' "s|__HARNESS_ROOT__|$PWD|g" .claude/scripts/auto_resume.plist

  launchd fallback that resumes the build after a usage-limit reset, even if
  the Claude session was hard-killed. Ticks every 30 minutes; auto_resume.sh
  no-ops when work is done, a resume is already running, or quota is still out.

  Enable (once):
    cp .claude/scripts/auto_resume.plist ~/Library/LaunchAgents/com.tddharness.autoresume.plist
    launchctl load ~/Library/LaunchAgents/com.tddharness.autoresume.plist
  Disable:
    launchctl unload ~/Library/LaunchAgents/com.tddharness.autoresume.plist
  Watch:
    tail -f .claude/state/auto_resume.log
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tddharness.autoresume</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__HARNESS_ROOT__/.claude/scripts/auto_resume.sh</string>
    </array>

    <!-- Every 30 minutes. The script itself is the real guard: it exits fast
         when there is nothing to do or quota is still out, so a frequent tick
         is cheap and means resume lands within ~30 min of the window reopening. -->
    <key>StartInterval</key>
    <integer>1800</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>__HARNESS_ROOT__/.claude/state/auto_resume.launchd.log</string>
    <key>StandardErrorPath</key>
    <string>__HARNESS_ROOT__/.claude/state/auto_resume.launchd.log</string>
</dict>
</plist>
